# docs/drive-notes-sync.timer
[Unit]
Description=Run Aletheia Google Drive case-notes sync every ~30s

[Timer]
# Fire 30s after boot, then every 30s after each run completes.
OnBootSec=30s
OnUnitActiveSec=30s
AccuracySec=5s
Persistent=true
Unit=drive-notes-sync.service

[Install]
WantedBy=timers.target
