# docs/qcluster.service
[Unit]
Description=Django-Q Cluster for Aletheia
After=network.target postgresql.service

[Service]
Type=simple
User=<user>
Group=<user>
WorkingDirectory=/home/<user>/law
ExecStart=/home/<user>/.venvs/law/bin/python manage.py qcluster
Restart=always
RestartSec=10

# Environment variables (adjust path if using .env differently)
EnvironmentFile=/home/<user>/law/config/.env

[Install]
WantedBy=multi-user.target
